# edgepulse-api.service
[Unit]
Description=EdgePulse FastAPI Backend
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/EdgePulse
EnvironmentFile=/root/EdgePulse/.env
Environment=PYTHONUNBUFFERED=1

ExecStart=/root/EdgePulse/.venv/bin/uvicorn api.main:app \
    --host 0.0.0.0 \
    --port 8502 \
    --workers 2 \
    --log-level info \
    --access-log \
    --proxy-headers \
    --forwarded-allow-ips=*

Restart=always
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=5
MemoryMax=800M
MemoryHigh=600M

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/root/EdgePulse
SyslogIdentifier=ep-api

[Install]
WantedBy=multi-user.target
