# edgepulse-intel.service
[Unit]
Description=EdgePulse Intel Node (signal generation)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/EdgePulse
EnvironmentFile=/root/EdgePulse/.env
Environment=PYTHONUNBUFFERED=1

ExecStart=/root/EdgePulse/.venv/bin/python3 -u ep_intel.py

Restart=always
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

# Resource limits
MemoryMax=1500M
MemoryHigh=1200M
CPUQuota=150%

# Watchdog — systemd kills and restarts if process hangs (no heartbeat for 5min)
WatchdogSec=300

# Logging
StandardOutput=append:/root/EdgePulse/output/logs/intel.log
StandardError=append:/root/EdgePulse/output/logs/intel.log
SyslogIdentifier=ep-intel

# Security
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
