# systemd/edgepulse-exec.service
[Unit]
Description=EdgePulse Exec — order execution
PartOf=edgepulse.target
After=docker.service network-online.target edgepulse-intel.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/EdgePulse
EnvironmentFile=/etc/edgepulse/edgepulse.env
Environment=PYTHONUNBUFFERED=1
Environment=MODE=exec
Environment=NODE_ID=exec-qvps-chi

# Wait for Redis socket before starting (max 60s)
ExecStartPre=/bin/bash -c \
    'for i in $(seq 1 30); do \
         redis-cli -s /var/run/redis/redis.sock ping 2>/dev/null | grep -q PONG && exit 0; \
         sleep 2; \
     done; echo "Redis not ready after 60s"; exit 1'

ExecStart=/root/EdgePulse/.venv/bin/python3 -u ep_exec.py

Restart=always
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

# Higher CPU priority than Intel — order placement is latency-sensitive
MemoryMax=512M
MemoryHigh=400M
CPUQuota=200%
Nice=-5

WatchdogSec=300

StandardOutput=append:/var/log/edgepulse/exec.log
StandardError=append:/var/log/edgepulse/exec.log
SyslogIdentifier=ep-exec

NoNewPrivileges=yes

[Install]
WantedBy=edgepulse.target
