# systemd/edgepulse-ob-depth.service
[Unit]
Description=EdgePulse OB Depth — order book depth imbalance KXFED/KXCPI/KXGDP
PartOf=edgepulse.target
After=docker.service network-online.target edgepulse-exec.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/EdgePulse
EnvironmentFile=/etc/edgepulse/edgepulse.env
Environment=PYTHONUNBUFFERED=1

ExecStart=/root/EdgePulse/.venv/bin/python3 -u ep_ob_depth.py

Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5
TimeoutStopSec=15

MemoryMax=512M
CPUQuota=100%
Nice=0

StandardOutput=append:/var/log/edgepulse/ob_depth.log
StandardError=append:/var/log/edgepulse/ob_depth.log
SyslogIdentifier=ep-ob-depth

NoNewPrivileges=yes

[Install]
WantedBy=edgepulse.target
