# systemd/edgepulse.target
[Unit]
Description=EdgePulse Trading System
Documentation=https://github.com/rbacidtech/Kalshi_bot
Wants=edgepulse-intel.service edgepulse-exec.service edgepulse-api.service edgepulse-llm.timer

[Install]
WantedBy=multi-user.target
